# deploy/smartposture-guardian.service
[Unit]
Description=SmartPosture Guardian - 慧姿·智能坐姿守护系统
Documentation=https://github.com/wyq-luminous/Group11
After=network-online.target hermes-gateway.service
Wants=network-online.target hermes-gateway.service
StartLimitIntervalSec=0

[Service]
Type=simple
# 确保 MCU 固件先启动（LED 点阵 + 蜂鸣器依赖 STM32 Bridge RPC）
# 前缀 '-': 忽略非零退出码。app 已在运行时 `app start` 会返回
# "App Is Running" 错误码 1，这属于期望状态，不应阻止服务启动。
ExecStartPre=-/usr/bin/arduino-app-cli app start user:posture_alerter
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/python3 /home/arduino/SmartPosture_Guardian/src/main.py
WorkingDirectory=/home/arduino/SmartPosture_Guardian
Environment="PATH=/home/arduino/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=/home/arduino"
Restart=always
RestartSec=5
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
